# contrib/templates/api.conf
# Starting conditions
start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

chdir /var/api
setuid transit

script
exec ./transit
end script

